# crontab
30 */2 * * * root /scripts/listener-quiz-clean.js >> /var/log/cron.log 2>&1
